# data/uos-installer-start.service
[Unit]
Description=launch installer

[Service]
Type=forking
ExecStart=/bin/bash /usr/bin/installer-start
User=root
[Install]
WantedBy=basic.target
